# cookbooks/cdo-apps/templates/default/puma.service.erb
[Unit]
Description=<%= @app_name %> service
# Ensure app-services are running before cloud-final user script
# reports SUCCESS to deployment process.
Before=cloud-final.service

# Configuration options copied from file generated by systemd-sysv-generator
[Service]
Type=forking
Restart=no
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=yes
ExecStart=/etc/init.d/<%= @app_name %> start
ExecStop=/etc/init.d/<%= @app_name %> stop

[Install]
WantedBy=multi-user.target